# etc/profile.d/ant.sh.symlink
ANT_HOME=/opt/ant
if [ -d $ANT_HOME/bin ] 
then
	PATH=$ANT_HOME/bin:$PATH
	export PATH
fi

